# roles/gitops_agent_gitlab/files/etc/systemd/system/gitops_agent.service
[Unit]
Description=GitOps Agent Service
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/etc/ansible/gitops_agent.sh