# roles/gitops_agent_gitlab/files/etc/systemd/system/gitops_agent.timer
[Unit]
Description=Hourly execute GitOps Agent from Genesis

[Timer]
OnUnitActiveSec=30m
Persistent=true

[Install]
WantedBy=timers.target